# etc/profile.d/flex-sdk.sh
export PATH=$PATH:~/Applications/flex_sdk_4.0.0.14159_mpl/bin
