# etc/profile.d/git-completion.sh
source $GIT_ROOT/contrib/completion/git-completion.bash

# want to have the branch name in the PS1 prompt?
#
#PS1='\h:\W$(__git_ps1 "(%s)") \u\$ '
