# etc/profile.d/postgresql-pgdata.sh
export PGDATA=$HOMEBREW_PREFIX/var/postgres